# akdev.autopatch/files/upgrade.cron
# Upgrade packages
0 0 * * * /usr/local/bin/yum-upgrade
