# canzuk/files/systemd/btrfs-snapshot@.service
[Unit]
Description= BTRFS Snapshot for %I

[Service]
ExecStart=/usr/local/bin/btrfs-snapshot --instance %I
